# simulador/crontab
* * * * * python3 /job/job_bank.py > /data/payment_records.json 2>&1
